# test/vm/resources/etc/profile.d/runc.sh
#!/bin/sh
alias runc="runc --root /run/containerd/runc/k8s.io/"
